# 10-htcondor-ce.conf
[program:htcondor-ce]
command=/usr/share/condor-ce/condor_ce_startup
autorestart=true
startsecs=20

